# templates/xvfb.service.j2
[Unit]
Description=Virtual Frame Buffer X Server
After=network.target

[Service]
ExecStart={{ xvfb_binary_path }} {{ xvfb_options }}

[Install]
WantedBy=multi-user.target
